# buildroot/board/custom_rpi/overlay/etc/init.d/S50read_sensors_web
#!/bin/sh
### BEGIN INIT INFO
# Provides:          read_sensors_web
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start read_sensors_web service at boot time
# Description:       Start the Python web application.
### END INIT INFO

NAME="read_sensors_web"
SCRIPT="/opt/read_sensors_web.py"
PIDFILE="/var/run/${NAME}.pid"
LOGFILE="/var/log/${NAME}.log"

start() {
    echo "Starting ${NAME}..."
    if [ -f ${PIDFILE} ] && kill -0 $(cat ${PIDFILE}) 2>/dev/null; then
        echo "${NAME} is already running."
        exit 1
    fi
    python3 /opt/read_sensors_job.py
    nohup python3 ${SCRIPT} >> ${LOGFILE} 2>&1 &
    echo $! > ${PIDFILE}
    echo "${NAME} started with PID $(cat ${PIDFILE})."
}

stop() {
    echo "Stopping ${NAME}..."
    if [ ! -f ${PIDFILE} ] || ! kill -0 $(cat ${PIDFILE}) 2>/dev/null; then
        echo "${NAME} is not running."
        exit 1
    fi
    kill $(cat ${PIDFILE})
    rm -f ${PIDFILE}
    echo "${NAME} stopped."
}

status() {
    if [ -f ${PIDFILE} ] && kill -0 $(cat ${PIDFILE}) 2>/dev/null; then
        echo "${NAME} is running with PID $(cat ${PIDFILE})."
    else
        echo "${NAME} is not running."
    fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 1
    start
    ;;
  status)
    status
    ;;
  *)
    echo "Usage: $0 {start|stop|restart|status}"
    exit 1
    ;;
esac

exit 0
