# buildroot/board/custom_rpi/overlay/etc/cron/crontabs/root
* * * * * echo "Kernel log message at Thu Jan  1 00:05:36 UTC 1970" > /dev/kmsg
