# buildroot/board/custom_rpi/overlay/etc/cron/crontabs/root
* * * * * python3 /opt/read_sensors_job.py
